# deploy/cnc-brain.service
[Unit]
Description=CNC AI Brain v11.0.1 - Union YouNi
After=network.target ollama.service

[Service]
Type=simple
User=timo
WorkingDirectory=/home/timo/projects/cnc-ai-brain
Environment=PATH=/home/timo/.miniconda/envs/forge_v7/bin:/usr/bin:/bin
ExecStart=/home/timo/.miniconda/envs/forge_v7/bin/python3 /home/timo/projects/cnc-ai-brain/app/main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
